# deploy/zoro2026-agent.service
[Unit]
Description=Zoro 2026 Raspberry Pi Agent
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/home/pi/zoro2026
ExecStart=/home/pi/zoro2026/.venv/bin/uvicorn pi_agent.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3
User=pi
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
