# deploy/zoro2026-body.service
[Unit]
Description=Zoro 2026 Pi Sensor Body Node
After=network-online.target zoro2026-agent.service
Wants=network-online.target

[Service]
WorkingDirectory=/home/pi/zoro2026
ExecStart=/home/pi/zoro2026/.venv/bin/python -m pi_agent.body_node
Restart=always
RestartSec=3
User=pi
Environment=PYTHONUNBUFFERED=1
Environment=LAPTOP_WS_BASE=ws://YOUR_LAPTOP_IP:8001
Environment=ZORO_AUDIO_SAMPLE_RATE=16000

[Install]
WantedBy=multi-user.target
